# systemd/dot-config/systemd/user/auto-theme-switcher.timer
[Unit]
Description=Auto adjusts GNOME theme between dark and light to match the time

[Timer]
OnCalendar=*-*-* 17:30:00
OnCalendar=*-*-* 09:00:00

[Install]
WantedBy=timers.target
